# auto-install/supervisor/control.conf
[program:control]
command=/usr/bin/python3 /usr/local/bin/pifire/control.py
directory=/usr/local/bin/pifire
autostart=true
autorestart=true
startretries=3
stderr_logfile=/usr/local/bin/pifire/logs/control.err.log
stdout_logfile=/usr/local/bin/pifire/logs/control.out.log
user=root
